# nodeinit.conf
description "node.js server"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

start script
   # script to create systemtray app will go here 
end script
